# deploy/pixelrag-api-green.service
[Unit]
Description=PixelRAG Search API (green slot — LoRA model)
After=network.target

[Service]
Type=simple
User=yichuan
WorkingDirectory=/home/yichuan/visrag
ExecStart=/home/yichuan/visrag/.venv/bin/pixelrag serve \
    --index-dir /home/yichuan/visrag-data/search_index_lora_vit_ckpt200_v2 \
    --tiles-dir /home/yichuan/visrag-data/tiles \
    --articles-json /home/yichuan/visrag-data/articles.json \
    --model /home/yichuan/visrag-data/models/Qwen3-VL-Embedding-2B-lora-vit-ckpt200 \
    --device cpu \
    --port 30002
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
